# docker/apache/src/supervisor/conf.d/apache2.conf
[program:apache2]
command=/usr/local/bin/apache2-foreground
stderr_logfile=/var/log/apache/apache/error.log
stdout_logfile=/var/log/apache/apache/access.log
